# src/911_reminder.service
[Service]
Restart=always
ExecStart=/usr/bin/python3 /home/dev/911_reminder/bot.py